# package/src/templates/backend-template/services/worker-service/crontab
# Graphile Worker Crontab
# Format: minute hour day month weekday task ?options {payload}
# All times are in UTC

# Production: Run task every hour at minute 0
# Uses fill=1d to backfill up to 1 day of missed runs if worker was down
0 * * * * taskName ?fill=1d

# Testing: Run task every 1 minute (commented out for production)
# * * * * * taskName ?fill=5m

# Alternative: Run task every 4 hours (more conservative)
# 0 */4 * * * taskName ?fill=1d